# tsv6.service
[Unit]
Description=TSV6 Raspberry Pi Video Player
Documentation=file:///home/%i/tsmedia-with-player/README.md
After=NetworkManager.service time-sync.target tsv6-boot-config-guard@%i.service
Wants=NetworkManager.service
Requires=tsv6-boot-config-guard@%i.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/tsmedia-with-player

# Environment
Environment="PATH=/home/%i/.local/bin:/home/%i/.cargo/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/%i/.Xauthority"
Environment="TSV6_ENVIRONMENT=production"
Environment="TSV6_SERVICE_NAME=tsv6.service"
Environment="SDL_AUDIODRIVER=dummy"
Environment="PYGAME_HIDE_SUPPORT_PROMPT=1"
Environment="ALSA_PCM_CARD=0"
Environment="ALSA_PCM_DEVICE=0"
Environment="ALSA_CONFIG_PATH=/dev/null"
Environment="VLC_ALSA_DEFAULT_DEVICE=/dev/null"
Environment="PULSE_SERVER="

# Switch network adapter (runs as root via sudo, non-fatal)
# Pass all connectivity env vars explicitly — sudo strips inherited env.
ExecStartPre=-/usr/bin/sudo TSV6_NETWORK_ADAPTER=${TSV6_NETWORK_ADAPTER} TSV6_WIFI_ENABLED=${TSV6_WIFI_ENABLED} TSV6_LTE_ENABLED=${TSV6_LTE_ENABLED} TSV6_CONNECTIVITY_MODE=${TSV6_CONNECTIVITY_MODE} /bin/bash /home/%i/tsmedia-with-player/scripts/switch-network-adapter.sh

# WIFI/LTE GATE: Block until network connectivity is confirmed.
# Sources /run/tsv6-network-adapter.env (written by switch-network-adapter.sh above).
# If WiFi-primary: waits for wlan0, starts provisioning if no saved WiFi.
# If LTE-primary: waits for wwan0, skips WiFi provisioning.
# Main app NEVER starts without network. (runs as root for systemctl access)
ExecStartPre=/usr/bin/sudo /bin/bash /home/%i/tsmedia-with-player/scripts/wifi-wait.sh

# CRITICAL: Wait for X11 socket to exist (up to 60 seconds)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 60); do [ -S /tmp/.X11-unix/X0 ] && exit 0; sleep 1; done; echo "X11 socket not found after 60s"; exit 1'

# Wait for X11 to be responsive (up to 30 seconds)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do xdpyinfo -display :0 >/dev/null 2>&1 && exit 0; sleep 1; done; echo "X11 not responsive after 30s"; exit 1'

# Disable screen blanking/DPMS before starting app (non-fatal)
ExecStartPre=-/bin/bash -c 'DISPLAY=:0 xset s off -dpms 2>/dev/null || true'

# Start application
ExecStart=/home/%i/.local/bin/uv run python run_production.py

# Restart policy - only restart on failure, not on clean exit
Restart=on-failure
RestartSec=10

# Allow wifi-wait.sh to block indefinitely during provisioning
TimeoutStartSec=infinity

# Watchdog - disabled (application handles its own monitoring via error recovery)
WatchdogSec=0
RuntimeMaxSec=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsv6

# Resource limits
TasksMax=300
MemoryMax=2G
CPUQuota=150%

# ── Connectivity Settings ─────────────────────────────────────────────────────
# Toggle WiFi and LTE independently. Both can be installed; only enabled ones
# are used at runtime. Change these to reconfigure, then:
#   sudo systemctl daemon-reload && sudo systemctl restart tsv6@<user>.service
#
# TSV6_WIFI_ENABLED   — Enable onboard WiFi (rpi-wifi)
# TSV6_LTE_ENABLED    — Enable 4G LTE modem (SIM7600NA-H / Hologram.io)
# TSV6_NETWORK_ADAPTER — Active adapter: rpi-wifi | intel-ax210 | 4g-lte
# TSV6_CONNECTIVITY_MODE — Failover strategy:
#   wifi_only                  WiFi only (LTE ignored even if enabled)
#   lte_only                   LTE only (WiFi ignored even if enabled)
#   wifi_primary_lte_backup    WiFi primary, fail over to LTE
#   lte_primary_wifi_backup    LTE primary, fail over to WiFi

Environment="TSV6_WIFI_ENABLED=true"
Environment="TSV6_LTE_ENABLED=false"
Environment="TSV6_NETWORK_ADAPTER=rpi-wifi"
Environment="TSV6_CONNECTIVITY_MODE=wifi_only"

# Derived network settings (written by switch-network-adapter.sh, dash means optional/missing OK)
EnvironmentFile=-/run/tsv6-network-adapter.env

# LTE Configuration (Hologram.io) — only used when TSV6_LTE_ENABLED=true
Environment="TSV6_LTE_APN=hologram"

# Bin level ToF sensor (M5Stack Mini ToF / VL53L0X on I2C bus 1)
Environment="TSV6_TOF_ENABLED=true"

# Recycle verification ToF sensor (M5Stack ToF4M / VL53L1X on I2C bus 2)
Environment="TSV6_RECYCLE_SENSOR_I2C_BUS=2"
Environment="TSV6_RECYCLE_SENSOR_THRESHOLD_MM=110"
Environment="TSV6_RECYCLE_SENSOR_POLL_INTERVAL=0.05"
Environment="TSV6_RECYCLE_SENSOR_SIMULATION=false"
Environment="TSV6_RECYCLE_SENSOR_REQUIRED_DETECTIONS=2"

# ── PiSignage Display Settings ───────────────────────────────────────────────
# Set PISIGNAGE_ENABLED=true to use a PiSignage-compatible display backend
# instead of local VLC+tkinter.
#
# PISIGNAGE_BACKEND selects the backend implementation:
#   rest   — PiSignageAdapter (REST client; requires a licensed PiSignage player
#             already registered with the Hostinger server)
#   native — TSV6NativeBackend (in-process Chromium+VLC player; no PiSignage
#             player license needed; recommended for new deployments)
#   vlc    — Legacy EnhancedVideoPlayer (tkinter+VLC; PISIGNAGE_ENABLED ignored)
#
# When PISIGNAGE_ENABLED=false (default), the legacy VLC player is used
# regardless of PISIGNAGE_BACKEND.
Environment="PISIGNAGE_ENABLED=true"
Environment="PISIGNAGE_BACKEND=native"
Environment="PISIGNAGE_SERVER_URL=https://tsmedia.g1tech.cloud"
Environment="PISIGNAGE_USERNAME=pi"
Environment="PISIGNAGE_PASSWORD=pi"
Environment="PISIGNAGE_INSTALLATION=g1tech26"
Environment="PISIGNAGE_GROUP=First Group"

# ── TSV6 Native Backend Settings ─────────────────────────────────────────────
# Used when PISIGNAGE_BACKEND=native.
# TSV6_VENUE_ID: Operator-assigned venue tag for impression reporting (optional)
# TSV6_APP_VERSION: Application version string injected into impression records
Environment="TSV6_VENUE_ID="
Environment="TSV6_APP_VERSION=1.0.0"

[Install]
WantedBy=multi-user.target
